# files/x86-64/files/etc/crontabs/root
5 17 * * * /usr/local/bin/singctl sb update >> /var/log/singctl_update.log 2>&1